# home/.config/systemd/user/clash.service
[Unit]
Description=fuckGFW

[Service]
Type=simple
ExecStart=/run/current-system/sw/bin/clash-meta

[Install]
WantedBy=default.target

